# installs/debian/itmon/DEBIAN/postinst
#!/bin/bash
set -e

# Enable and start the service
systemctl enable itmon.service
systemctl start itmon.service

cp ../../../../main.py /usr/lib/itmon.reporter.py
exit 0