# installs/debian/itmon/DEBIAN/prerm
# Stop the service
systemctl stop itmon

# Remove the service
systemctl disable itmon
rm /usr/lib/systemd/system/itmon.service